# crontab
*/55 * * * * cd /app && uv run home-finder 2>&1
